# services/geolocation.timer
[Unit]
Description=Run MTK Tankstellen Geocoding every hour

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
